# scripts/com.bridge.macbridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bridge.macbridge</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>USER_HOME/.mac-bridge/mac-bridge.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>USER_HOME/.mac-bridge/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>USER_HOME/.mac-bridge/launchd_stderr.log</string>
</dict>
</plist>
